# etc/supervisor/conf.d/timelineserver.conf
[program:timelineserver]
command=/opt/hadoop-2.7.0/bin/yarn timelineserver